# pm2.conf
[supervisord]
nodaemon=true

[program:pm2]
command:pm2 start /app/backend/bin/www --no-daemon
startretries:5
